# root/etc/cron.d/flashstart-update-ip
*/30 * * * * root /etc/e-smith/events/actions/nethserver-flashstart-update-ip
